# scripts/service/xrdp-sesman/run
#!/bin/sh

set -eu

exec 2>&1
exec xrdp-sesman --config /opt/xrdp/etc/xrdp/sesman.ini --nodaemon
