# scripts/service/xrdp-sesman/run
#!/bin/sh

exec /usr/sbin/xrdp-sesman --nodaemon
