# scripts/service/udevadm-trigger/run
#!/bin/sh

sv start udevd || exit 1

udevadm trigger
exec sleep infinity
