# scripts/service/udevadm-trigger/run
#!/bin/sh

set -eu

sv start udevd >/dev/null

exec 2>&1
udevadm trigger ||:
exec chpst -b udevadm-trigger perl -MPOSIX -e 'pause()'
