# scripts/service/udevadm-trigger/run
#!/bin/sh

sv start udevd >/dev/null || exit 1

/usr/bin/udevadm trigger
exec sleep infinity
