# scripts/service/dbus-daemon/run
#!/bin/sh

exec /usr/bin/sudo -u messagebus /usr/bin/dbus-daemon --system --nofork
