# scripts/service/sshd/run
#!/bin/sh

set -eu

if [ ! -d /run/sshd/ ]; then
	install -m 755 -o root -g root -d /run/sshd/
fi

exec 2>&1
exec /usr/sbin/sshd -D
