# scripts/service/xrdp/run
#!/bin/sh

set -eu

exec 2>&1
exec xrdp --config /opt/xrdp/etc/xrdp/xrdp.ini --nodaemon
