# scripts/service/xrdp/run
#!/bin/sh

exec /usr/sbin/xrdp --nodaemon
